# build-files/conf/supervisor.d/crond.conf
[program:cron]
autorestart=%(SUPERVISOR_AUTOSTART_CRON)s
command=cron -f